# bismark-mserver/etc/cron.d/bismark-mserver
# m  h  dom mon dow   command
*/5  *  *   *   *     bismark-mserver run > /dev/null 2>&1 &
*/10 *  *   *   *     bismark-mserver-heartbeat > /dev/null 2>&1 &
