# bismark-mserver/etc/cron.d/bismark-mserver
# m  h    dom mon dow   user command
*/5  *    *   *   *     root bismark-mserver run > /dev/null 2>&1 &
*/5  *    *   *   *     root bismark-mserver-hostsallow > /dev/null 2>&1 &
*/2  *    *   *   *     root bismark-mserver-heartbeat > /dev/null 2>&1 &
#0    */6  *   *   *     root bismark-mserver-perfcheck > /dev/null 2>&1 &
